# crontab
# Update domains
*/5 * * * *   /usr/bin/freedns.sh